# config/svcs-solaris/lwsmd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Copyright 2009 Likewise Software, All rights reserved.
        Use is subject to license terms.

        #ident  "@(#)lwsmd.xml        5.3     08/09/21 SMI"
-->

<service_bundle type='manifest' name='LIKEWISElwis:lwsmd'>

<service name='network/lwsmd' type='service' version='1'>

        <create_default_instance enabled='false' />

        <single_instance/>

        <dependency name='net-loopback' grouping='require_any'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/loopback' />
        </dependency>

        <dependency name='net-service' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/service'/>
        </dependency>

        <dependency name='net-physical' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/physical' />
        </dependency>

        <dependency name='filesystem-local' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/system/filesystem/local' />
        </dependency>

        <dependent name='lwsmd-multi-user-server' grouping='optional_all'
                        restart_on='none'>
                <service_fmri value='svc:/milestone/multi-user-server' />
        </dependent>

        <exec_method type='method' name='start'
                exec='/opt/likewise/sbin/lwsmd --start-as-daemon'
                timeout_seconds='60'>
                <method_context>
                        <method_environment>
                                <envvar name="LD_LIBRARY_PATH" value=""/>
                                <envvar name="LD_PRELOAD" value=""/>
                        </method_environment>
                </method_context>
        </exec_method>

        <exec_method type='method' name='stop'
                exec='/opt/likewise/bin/lwsm shutdown'
                timeout_seconds='60'>
                <method_context>
                        <method_environment>
                                <envvar name="LD_LIBRARY_PATH" value=""/>
                                <envvar name="LD_PRELOAD" value=""/>
                        </method_environment>
                </method_context>
        </exec_method>

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                        Likewise Software Service Management daemon
                        </loctext>
                </common_name>
                <documentation>
                        <doc_link name='Likewise Software Documentation' uri='http://www.likewise.com/resources/documentation_library' />
                </documentation>
        </template>
</service>

</service_bundle>
