# package/open/deb/pbis-open.postrm.in
#!/bin/sh

UPGRADEDIR=/var/lib/pbis-upgrade

LOG=/var/log/pbis-open-install.log
TLOG=/tmp/pbis-open-install.tmp

PKG_ARCH="@MK_HOST_ARCH@"

DAEMONS_TO_HALT="lwsmd lwregd netlogond lwiod dcerpcd eventlogd lsassd srvsvcd"

# Display to screen and log file with a blank line between entries.
log()
{
    echo $@
    echo
    echo $@ >> $LOG
    echo >> $LOG
}

# Display to screen and log file with no blank line.
_log()
{
    echo $@
    echo $@ >> $LOG
}

# Display to file.
logfile()
{
    echo $@ >> $LOG
    echo >> $LOG
}

# Execute command.
# If successful, note in log file.
# If not successful, note on screen and log file.
run()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
        cat $TLOG >> $LOG
        echo >> $LOG
    else
        _log "Error: $@ returned $err"
        _log `cat $TLOG`
        _log
    fi
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

# Execute command.
# Log only to file.
run_quiet()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
    else
        echo "Error: $@ returned $err  (ignoring and continuing)" >> $LOG
    fi
    cat $TLOG >> $LOG
    echo >> $LOG
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

# Execute command.
# If successful, note in log file.
# If not successful, note on screen and log file and then exit.
run_or_fail()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
        cat $TLOG >> $LOG
        echo >> $LOG
    else
        _log "Error: $@ returned $err  (aborting this script)"
        _log `cat $TLOG`
        _log
        rm -f $TLOG > /dev/null 2>&1
        exit 1
    fi
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

case "$1" in
    remove)
        logfile "Package: PowerBroker Identity Services Open postremove [remove] begins (`date`)"

        run_quiet /usr/sbin/service pbis stop
        run_quiet /usr/sbin/service lwsmd stop

        for daemon in $DAEMONS_TO_HALT
        do
            run_quiet pkill -KILL -x $daemon
        done

        run_quiet /usr/sbin/update-rc.d pbis remove
        run_quiet /usr/sbin/update-rc.d lwsmd remove

        logfile "Package: PowerBroker Identity Services Open postremove [remove] finished"
    ;;

    purge)
        logfile "Package: PowerBroker Identity Services Open postremove [purge] begins (`date`)"

        run_quiet /usr/sbin/service pbis stop
        run_quiet /usr/sbin/service lwsmd stop

        for daemon in $DAEMONS_TO_HALT
        do
            run_quiet pkill -KILL -x $daemon
        done

        run_quiet /usr/sbin/update-rc.d pbis remove
        run_quiet /usr/sbin/update-rc.d lwsmd remove

        run_quiet rm -rf '@LW_CACHEDIR@'

        logfile "Package: PowerBroker IdentityServices Open postremove [purge] finished"
    ;;
esac

#DEBHELPER#
