# package/open/deb/pbis-open.prerm.in
#!/bin/sh

UPGRADEDIR=/var/lib/pbis-upgrade

LOG=/var/log/pbis-open-install.log
TLOG=/tmp/pbis-open-install.tmp

PKG_ARCH="@MK_HOST_ARCH@"

DAEMONS_TO_HALT="reapsysld lsassd lwiod netlogond eventlogd dcerpcd lwregd lwsmd"

# Display to screen and log file with a blank line between entries.
log()
{
    echo $@
    echo
    echo $@ >> $LOG
    echo >> $LOG
}

# Display to screen and log file with no blank line.
_log()
{
    echo $@
    echo $@ >> $LOG
}

# Display to file.
logfile()
{
    echo $@ >> $LOG
    echo >> $LOG
}

# Execute command.
# If successful, note in log file.
# If not successful, note on screen and log file.
run()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
        cat $TLOG >> $LOG
        echo >> $LOG
    else
        _log "Error: $@ returned $err"
        _log `cat $TLOG`
        _log
    fi
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

# Execute command.
# Log only to file.
run_quiet()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
    else
        echo "Error: $@ returned $err  (ignoring and continuing)" >> $LOG
    fi
    cat $TLOG >> $LOG
    echo >> $LOG
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

# Execute command.
# If successful, note in log file.
# If not successful, note on screen and log file and then exit.
run_or_fail()
{
    "$@" > $TLOG 2>&1
    err=$?
    if [ $err -eq 0 ]; then
        echo "Success: $@" >> $LOG
        cat $TLOG >> $LOG
        echo >> $LOG
    else
        _log "Error: $@ returned $err  (aborting this script)"
        _log `cat $TLOG`
        _log
        rm -f $TLOG > /dev/null 2>&1
        exit 1
    fi
    rm -f $TLOG > /dev/null 2>&1
    return $err
}

execute_auxiliary_scripts()
{
    if [ -d "@LW_CACHEDIR@/scripts/uninstall" ]; then
        for file in @LW_CACHEDIR@/scripts/uninstall/*; do
            run_quiet "$file" --uninstall
        done
    fi
}

prerm_remove()
{
    logfile "Package: PowerBroker Identity Services Open preremove [remove] begins (`date`)"

    execute_auxiliary_scripts

    run_quiet '@MK_BINDIR@/domainjoin-cli' configure --disable pam
    run_quiet '@MK_BINDIR@/domainjoin-cli' configure --disable nsswitch

    run_quiet /usr/sbin/service pbis stop
    run_quiet /usr/sbin/service lwsmd stop

    run_quiet /usr/sbin/update-rc.d pbis remove
    run_quiet /usr/sbin/update-rc.d lwsmd remove

    run_quiet rm -f '@LW_INITDIR@/pbis'
    run_quiet rm -f '@LW_INITDIR@/lwsmd'

    for daemon in $DAEMONS_TO_HALT
    do
        run_quiet pkill -KILL -x $daemon
    done

    logfile "Package: PowerBroker Identity Services Open preremove [remove] finished"
    exit 0
}

case "$1" in
    remove)
        prerm_remove
        ;;

    failed-upgrade)
    # The prerm script of a previous version failed for upgrade.
        ;;

    upgrade)
    # A newer version is about to be installed; do nothing and let it handle it.
        ;;
esac
