# config/linux/suse/pbis.in
#! /bin/sh
#
# Copyright (c) BeyondTrust Software.  All rights reserved.
#
### BEGIN INIT INFO
# Provides: pbis
# Required-Start: $network $remote_fs lwsmd
# Required-Stop:
# Default-Start: 3 5
# Default-Stop:
# Description: Start PowerBroker Identity Services Open services
### END INIT INFO

PREFIX="@MK_PREFIX@"
SERVICE_NAME="pbis"

. @MK_LIBEXECDIR@/init-pbis.sh
