# config/linux/redhat/lwsmd.service.in
[Unit]
Description=BeyondTrust PBIS Service Manager
After=network.target 
Before=sshd.service

[Service]
Type=forking
ExecStart=@MK_SBINDIR@/lwsmd --start-as-daemon 
ExecReload=@MK_BINDIR@/lwsm refresh
ExecStop=@MK_BINDIR@/lwsm shutdown
# We want systemd to give lwsmd some time to finish gracefully, but still want
# it to kill lwsmd after TimeoutStopSec if something went wrong during the
# graceful stop. Normally, Systemd sends SIGTERM signal right after the
# ExecStop, which would kill lwsmd. We are sending useless SIGCONT here to give
# lwsmd time to finish.
KillSignal=SIGCONT
PrivateTmp=false

[Install]
WantedBy=multi-user.target nss-lookup.target
