# config/linux/redhat/pbis.in
#! /bin/sh
#
# Copyright (c) BeyondTrust Software.  All rights reserved.
#
###
# chkconfig: 35 21 9
# description: Start PowerBroker Identity Services Open services
###
#

PREFIX="@MK_PREFIX@"
SERVICE_NAME="pbis"

. @MK_LIBEXECDIR@/init-pbis.sh
