# config/solaris/pbis.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Copyright 2009 Likewise Software, All rights reserved.
        Use is subject to license terms.

        #ident  "@(#)likewise.xml        5.3     08/09/21 SMI"
-->

<service_bundle type='manifest' name='BEYONDTRUSTpbis:pbis'>

<service name='network/pbis' type='service' version='1'>

        <create_default_instance enabled='false' />

        <single_instance/>

        <dependency name='pbis-lwsmd' grouping='require_any'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/lwsmd' />
        </dependency>

        <dependency name='net-loopback' grouping='require_any'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/loopback' />
        </dependency>

        <dependency name='net-service' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/service'/>
        </dependency>

        <dependency name='net-physical' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/network/physical' />
        </dependency>

        <dependency name='filesystem-local' grouping='require_all'
                        restart_on='none' type='service'>
                <service_fmri value='svc:/system/filesystem/local' />
        </dependency>

        <dependent name='likewise-multi-user-server' grouping='optional_all'
                        restart_on='none'>
                <service_fmri value='svc:/milestone/multi-user-server' />
        </dependent>

        <exec_method type='method' name='start'
                exec='/opt/pbis/bin/lwsm -q autostart'
                timeout_seconds='60' />

        <exec_method type='method' name='stop'
                exec=':true'
                timeout_seconds='60' />

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                        PowerBroker Identity Services autostart
                        </loctext>
                </common_name>
                <documentation>
                        <doc_link name='PowerBroker Identity Services Documentation' uri='http://www.beyondtrust.com/Products/PowerBroker-Identity-Services-Open-Edition/Resources/Documentation-Library/' />
                </documentation>
        </template>
</service>

</service_bundle>
