# git-build/package/open/deb/likewise-open.postinst.in
#!/bin/sh

#set -e

PATH=$PATH:@MK_BINDIR@
export PATH

case "$1" in
    abort-upgrade)
        ;;

    configure)

        /etc/init.d/lwsmd stop
        sleep 2

        if [ "@IS_EMBEDDED@" = "no" ]
        then
            domainjoin-cli configure --enable pam
            domainjoin-cli configure --enable nsswitch
        fi

        @MK_SBINDIR@/lwsmd --start-as-daemon || exit 1

        for file in "@LW_CONFIGDIR@/"*.reg; do
            echo "Importing $file..."
            lwregshell import $file || exit 1
        done

        lwsm shutdown || exit 1
        sleep 2

        /usr/sbin/update-rc.d -f lwsmd defaults 17 8
        /usr/sbin/update-rc.d -f likewise defaults 18 7

        /etc/init.d/lwsmd start
        /etc/init.d/likewise star
        ;;
esac

#DEBHELPER#
